# mcp-bridge.service
[Unit]
Description=MCP Bridge Server
After=network.target

[Service]
Type=simple
DynamicUser=yes
WorkingDirectory=/opt/mcp-bridge
ExecStart=/opt/mcp-bridge/venv/bin/python3 server.py
Restart=on-failure
RestartSec=5

# ping capability
AmbientCapabilities=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
