# mcp-bridge.service
[Unit]
Description=MCP Bridge Server
After=network.target

[Service]
Type=simple
User=nobody
WorkingDirectory=/opt/mcp-bridge
ExecStart=/opt/mcp-bridge/venv/bin/python3 server.py
Restart=on-failure
RestartSec=5

# ping capability를 nobody 사용자에게 상속
AmbientCapabilities=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
