# debian/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade)
        systemctl stop mcp-bridge.service 2>/dev/null || true
        systemctl disable mcp-bridge.service 2>/dev/null || true
        ;;
esac

exit 0
